# x2goserver.service
[Unit]
Description=X2Go Server Daemon

[Service]
ExecStart=/usr/sbin/x2gocleansessions

[Install]
WantedBy=multi-user.target

